# battleship-deluxe.service
[Unit]
Description=Battleship Deluxe (Node Web + WebSocket)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/battleship-deluxe
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/battleship-deluxe/.env
ExecStart=/usr/bin/node /var/www/battleship-deluxe/server.js
Restart=on-failure
RestartSec=3
User=www-data
Group=www-data

# Hardening (optional but nice)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/battleship-deluxe

[Install]
WantedBy=multi-user.target
